# systemd/nhd12864-demo.service
[Unit]
Description=NHD12864 LCD demo (C++ spidev + libgpiod)
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=/etc/nhd12864/nhd12864.conf
ExecStart=/opt/nhd12864/nhd12864_demo --spidev ${SPIDEV} --chip ${GPIOCHIP} --dc ${DC_LINE} --rst ${RST_LINE}
Restart=on-failure

[Install]
WantedBy=multi-user.target
